# deploy/systemd/corpsite-healthcheck.timer
[Unit]
Description=Run Corpsite health check every 5 minutes
Documentation=file:/opt/projects/corpsite/app/docs/deploy/VPS_STABILITY.md

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
