# deploy/systemd/corpsite-frontend.service
# Corpsite Next.js frontend (production).
#
# Install:
#   sudo cp deploy/systemd/corpsite-frontend.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable corpsite-frontend
#
# Deploy (always build before start):
#   git pull && sudo ./scripts/deploy_frontend.sh
#
# ADR-INFRA-005: port guard + KillMode prevent EADDRINUSE orphan loops.

[Unit]
Description=Corpsite Next.js Frontend
Documentation=file:/opt/projects/corpsite/app/docs/deploy/VPS_STABILITY.md
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/projects/corpsite/app/corpsite-ui
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=CORPSITE_UI_DIR=/opt/projects/corpsite/app/corpsite-ui

ExecStartPre=/opt/projects/corpsite/app/scripts/check_frontend_build.sh
ExecStartPre=+/opt/projects/corpsite/app/scripts/ops/ensure_port_free.sh 3000 --service corpsite-frontend --orphan-pattern next-server --orphan-pattern "npm run start" --orphan-pattern node
ExecStart=/usr/bin/npm run start --silent

Restart=on-failure
RestartSec=5
RestartPreventExitStatus=203
KillMode=mixed
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=corpsite-frontend

[Install]
WantedBy=multi-user.target
